# deploy/gpsa-api.service
[Unit]
Description=GPSA-UDS FastAPI backend
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/gpsa-web/gpsa
EnvironmentFile=/etc/gpsa/gpsa.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/home/ubuntu/gpsaenv/bin/alembic upgrade head
ExecStart=/home/ubuntu/gpsaenv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30
KillSignal=SIGTERM

# Service hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/projects/gpsa-web/gpsa/storage/temp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
